# packaging/coredns-ztnet.service
[Unit]
Description=CoreDNS with ztnet plugin (ZeroTier DNS)
Documentation=https://github.com/CleoWixom/ztnet-dns
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=coredns
Group=coredns
ExecStart=/usr/sbin/coredns -conf /etc/coredns/Corefile
ExecReload=/bin/kill -USR1 $MAINPID

Restart=on-failure
RestartSec=5s

# Logs: journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coredns-ztnet

# Hardening
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
# Required for CAP_NET_BIND_SERVICE (port 53), granted via setcap in postinstall
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
